# rootfs/etc/profile.d/geodesic.sh
#!/usr/bin/env bash
export GEODESIC_SHELL=true
complete -r -A file